# sebastian-api.service
[Unit]
Description=Sebastian Website API
After=network.target postgresql.service

[Service]
Type=simple
User=noetrevino
WorkingDirectory=/opt/sebastian/prod
EnvironmentFile=/opt/sebastian/prod/.env
ExecStart=/usr/bin/java -jar /opt/sebastian/prod/backend/build/libs/sebastian-api.jar
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
